# service/openclaw-claude-bridge.service
[Unit]
Description=OpenClaw Claude Bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/openclaw-claude-bridge
ExecStart=/root/.nvm/versions/node/v24.14.1/bin/node src/index.js
KillMode=process
TimeoutStopSec=180
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=HOME=/root
Environment=CLAUDE_BIN=/root/.nvm/versions/node/v24.14.1/bin/claude
Environment=DEBUG_IO=1

[Install]
WantedBy=multi-user.target
